# asr.service
[Unit]
Description=ASR Pipeline Service - Local Speech Recognition
Documentation=https://github.com/your-repo/asr-pipeline
After=network.target sound.service
Wants=network.target

[Service]
Type=simple
User=asr
Group=asr

# Working directory - adjust path as needed
WorkingDirectory=/opt/asr
ExecStart=/opt/asr/venv/bin/python /opt/asr/asr_service.py

# Environment variables
Environment=PYTHONPATH=/opt/asr
Environment=CUDA_VISIBLE_DEVICES=0
Environment=PYTHONUNBUFFERED=1

# Logging
StandardOutput=append:/opt/asr/logs/asr.out
StandardError=append:/opt/asr/logs/asr.err

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/asr/logs
ReadOnlyPaths=/opt/asr

# Resource limits
MemoryMax=8G
CPUQuota=200%

# Audio device access
SupplementaryGroups=audio

[Install]
WantedBy=multi-user.target
