# asr.service
[Unit]
Description=ASR Pipeline Service - Local Speech Recognition
Documentation=https://github.com/Sang-Buster/Aerolex-Logger
After=network.target sound.target
Wants=network.target

[Service]
Type=simple
User=@CURRENT_USER@
Group=@CURRENT_USER@

# Project directory (set by install script)
WorkingDirectory=@ASR_DIR@

# Run inside venv
ExecStart=@ASR_DIR@/.venv/bin/python @ASR_DIR@/src/asr_service.py

# Environment
Environment="PYTHONPATH=@ASR_DIR@"
Environment="CUDA_VISIBLE_DEVICES=0"
Environment="PYTHONUNBUFFERED=1"
Environment="PULSE_RUNTIME_PATH=/run/user/@USER_ID@/pulse"
Environment="XDG_RUNTIME_DIR=/run/user/@USER_ID@"
Environment="HOME=/home/@CURRENT_USER@"
Environment="USER=@CURRENT_USER@"
Environment="DISPLAY=:0"

# Logging (systemd journal only — safer)
StandardOutput=journal
StandardError=journal

# Restart policy
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
